# feeds_src/packages/ipv6/ipv6tunnel/files/etc/init.d/tunnelscript
#!/bin/sh /etc/rc.common
# Handles manual IPIP6 Tunneling
#

START=100
STOP=1

start() {
	/sbin/tunnellingscript start
}

stop() {
	/sbin/tunnellingscript stop
}    
